# base-image/resources/supervisor/programs/vncserver.conf
[program:vncserver]
command=python %(ENV_RESOURCES_PATH)s/scripts/vncrun.py
user=%(ENV_NB_USER)s
; user=root
priority=55
autostart=true
autorestart=true
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)
stopasgroup=true ; If true, the flag causes supervisor to send the stop signal to the whole process group and implies killasgroup is true. This is useful for programs, such as Flask in debug mode, that do not propagate stop signals to their children, leaving them orphaned.
killasgroup=true
;environment=HOME="/home/%(ENV_NB_USER)s",USER=%(ENV_NB_USER)s,XDG_RUNTIME_DIR="/tmp/xdg_runtime_dir"
